# ops/systemd/dos-amigos-daily-coo-reconcile.timer
[Unit]
Description=Run Dos Amigos daily COO morning reconciliation

[Timer]
OnCalendar=*-*-* 06:00:00 Asia/Yekaterinburg
Persistent=true
Unit=dos-amigos-daily-coo-reconcile.service

[Install]
WantedBy=timers.target
