# ops/systemd/dos-amigos-daily-coo-reconcile.service
[Unit]
Description=Dos Amigos daily COO morning reconciliation
Wants=network-online.target
After=network-online.target docker.service
StartLimitIntervalSec=3600
StartLimitBurst=5

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/dos-amigos-core
ExecStart=/opt/dos-amigos-core/ops/bin/run-daily-coo-report.sh --reconcile --json
Restart=on-failure
RestartSec=60
TimeoutStartSec=7200
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
