# ops/systemd/dos-amigos-sales-sync.service
[Unit]
Description=Dos Amigos one-shot iiko sales automation
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/dos-amigos-core
ExecStart=/opt/dos-amigos-core/ops/bin/run-sales-sync.sh
TimeoutStartSec=1800
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
